# libvirt-resolved-bridge.service
[Unit]
Description=Bridge Libvirt Network Config to Systemd-Resolved
Documentation=https://github.com/sjd78/libvirt-resolved-bridge
After=libvirtd.service libvirt-dbus.service systemd-resolved.service
Requires=libvirt-dbus.service systemd-resolved.service

[Service]
Type=simple
ExecStart=/usr/bin/libvirt-resolved-bridge
Restart=on-failure
RestartSec=5s

# Run as root with libvirt group for DBus policy
SupplementaryGroups=libvirt

# SELinux: run unconfined to allow DBus communication with libvirt
SELinuxContext=system_u:system_r:unconfined_t:s0

# Security hardening
NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
